# files/system/etc/systemd/system/speedtest.timer
[Unit]
Description=speedtest

[Timer]
OnUnitActiveSec=1h
OnBootSec=1h

[Install]
WantedBy=timers.target