# files/system/etc/systemd/system/speedtest.service
[Unit]
Description=Runs an internet speedtest

[Service]
Type=oneshot
ExecStart=/usr/bin/speedtest-cli